# systemd/rpicserver.service
[Unit]
Description=Rpic Server

[Service]
Type=simple
ExecStart=/usr/bin/rpic-server-linux-arm

[Install]
WantedBy=multi-user.target
